# linux/software/openfortivpn/templates/openfortivpn.service
[Unit]
Description=OpenFortiVpn

[Service]
ExecStart=/usr/bin/openfortivpn

[Install]
WantedBy=multi-user.target